# contrib/r-tftpd.socket
[Unit]
Description=r-tftpd Server Activation Socket

[Socket]
ListenDatagram=69

[Install]
WantedBy=sockets.target
