# contrib/r-tftpd.service.in
[Unit]
Description=r-tftpd Server
Requires=r-tftpd.socket

[Service]
Environment=RUST_LOG=info
ExecStart=/usr/sbin/r-tftpd --systemd $TFTPD_OPTS
WorkingDirectory=/var/lib/tftpboot
User=ftp
Group=ftp
